# CuulFan.service
[Unit]
Description=CuulFan

[Service]
ExecStart=/home/pi/CuulFan/CuulFan.py
Restart=always

[Install]
WantedBy=multi-user.target
